# integrations/capos/systemd/argosfs-health.service
[Unit]
Description=ArgosFS rootfs health check
After=local-fs.target
ConditionPathExists=/usr/sbin/argosfs

[Service]
Type=oneshot
Environment=ARGOSFS_POOL_CONFIG=/etc/argosfs/root-pool.json
Environment=ARGOSFS_BACKEND=raw
Environment=ARGOSFS_DEVICES=
ExecStart=/bin/sh -c 'if [ -r "$ARGOSFS_POOL_CONFIG" ]; then exec /usr/sbin/argosfs fsck --pool-config "$ARGOSFS_POOL_CONFIG"; elif [ -n "$ARGOSFS_DEVICES" ]; then exec /usr/sbin/argosfs fsck --backend "$ARGOSFS_BACKEND" --devices "$ARGOSFS_DEVICES"; else exit 0; fi'
