# integrations/capos/systemd/argosfs-watchdog.service
[Unit]
Description=ArgosFS rootfs watchdog
After=local-fs.target
ConditionPathExists=/usr/sbin/argosfs

[Service]
Type=simple
Environment=ARGOSFS_POOL_CONFIG=/etc/argosfs/root-pool.json
Environment=ARGOSFS_BACKEND=raw
Environment=ARGOSFS_DEVICES=
ExecStart=/bin/sh -c 'while sleep 60; do if [ -r "$ARGOSFS_POOL_CONFIG" ]; then /usr/sbin/argosfs preflight-root --pool-config "$ARGOSFS_POOL_CONFIG" --mode ro --json >/run/argosfs-watchdog.json || exit 1; elif [ -n "$ARGOSFS_DEVICES" ]; then /usr/sbin/argosfs preflight-root --backend "$ARGOSFS_BACKEND" --devices "$ARGOSFS_DEVICES" --mode ro --json >/run/argosfs-watchdog.json || exit 1; fi; done'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
