# contrib/capos/systemd/argosfs-watchdog.service
[Unit]
Description=ArgosFS rootfs watchdog
After=local-fs.target
ConditionPathExists=/usr/sbin/argosfs

[Service]
Type=simple
Environment=ARGOSFS_BACKEND=raw
Environment=ARGOSFS_DEVICES=
ExecStart=/bin/sh -c 'while sleep 60; do test -n "$ARGOSFS_DEVICES" && /usr/sbin/argosfs preflight-root --backend "$ARGOSFS_BACKEND" --devices "$ARGOSFS_DEVICES" --mode ro >/run/argosfs-watchdog.json || exit 1; done'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
