# contrib/capos/systemd/argosfs-recovery.target
[Unit]
Description=ArgosFS recovery mode
Requires=emergency.target
After=emergency.target
AllowIsolate=yes
